# etc/systemd/system/ping_watchdog.service
[Unit]
Description=Proxmox guest agent ping watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/ping_watchdog
ExecStart=/usr/local/bin/ping_watchdog.py
